# supervisord.ini
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/var/run/supervisord.pidfile

[program:webroot]
command=sh /usr/local/bin/entrypoint.sh
priority=1

[program:node]
command=npm run dev